# barman/crontab
* * * * * barman /usr/local/bin/barman -q cron
@daily root /usr/sbin/logrotate /etc/logrotate.conf
